# claribot.service.template
[Unit]
Description=Claribot - LLM Project Automation Service
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__HOME__/.claribot
ExecStart=/usr/local/bin/claribot
Restart=on-failure
RestartSec=5
Environment=HOME=__HOME__
Environment=PATH=__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
